# roles/sysmon_linux/templates/sysmon.service.j2
[Unit]
Description=Sysmon for Linux
Documentation=https://github.com/microsoft/SysmonForLinux
After=network.target syslog.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/sysmon
ExecReload=/usr/bin/sysmon -c {{ sysmon_config_path }}
Restart=on-failure
RestartSec=5

# Sysmon olaylarını journald yerine doğrudan dosyaya yaz.
# Bu sayede journald şişmez; rsyslog imfile ile okuyup SIEM'e iletir.
StandardOutput=append:{{ sysmon_log_path }}
# Servis hataları (crash, start fail) debug için journal'da kalsın.
# Ancak rate limit ile patlamayı önle: 30 saniyede max 50 satır.
StandardError=journal
SyslogIdentifier=sysmon
LogRateLimitIntervalSec={{ sysmon_journal_rate_interval_sec }}
LogRateLimitBurst={{ sysmon_journal_rate_burst }}

[Install]
WantedBy=multi-user.target